# storage/supervisord.conf
# =============================================================================
# ATLAS Analytics Container — Supervisord Configuration
# =============================================================================
# Manages four processes inside the unified analytics container:
#   1. PostgreSQL   — Metadata store (watermarks, device/location registry)
#   2. ClickHouse   — Analytical query engine
#   3. Delta Loader — Batch ingestion from Refined Layer into ClickHouse
#   4. ML Loader    — ML inference Parquet ingestion into ClickHouse
# =============================================================================

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:postgresql]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/data
user=postgres
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:clickhouse]
command=/usr/bin/clickhouse-server --config-file=/etc/clickhouse-server/config.xml
user=clickhouse
autorestart=true
startretries=30
startsecs=5
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:delta-loader]
command=/app/loader-start.sh
autorestart=true
priority=20
startsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ml-loader]
command=/app/ml-loader-start.sh
autorestart=true
priority=25
startsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:streamlit]
command=python3 -m streamlit run /app/app.py --server.port=8501 --server.address=0.0.0.0
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
